# piavpn/piavpn_custom.service
[Unit]
Description=Starts Pia VPN

[Service]
ExecStart=/etc/private-internet-access/piactl_startup.sh

[Install]
WantedBy=multi-user.target 
